# jmx_prometheus_httpserver/src/main/docker/etc/supervisor/conf.d/program-jmx-exporter.conf
[program:jmx-exporter]
command=/usr/local/bin/launch-jmx-exporter

redirect_stderr=true
stdout_logfile=/var/log/supervisor/jmx-exporter.log
stderr_logfile=NONE
